# systemd/odoo-admin-worker@.service
[Unit]
Description=Odoo SaaS Admin RQ Worker %i
After=network.target redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/odoo-admin
EnvironmentFile=/opt/odoo-admin/.env
ExecStart=/opt/odoo-admin/venv/bin/rq worker -u ${REDIS_URL} odoo_admin_jobs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
